# work/com.semiconductor.investor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.semiconductor.investor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhen/Documents/Codex/2026-06-07/deepseek/work/auto_runner.sh</string>
        <string>weekly</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <array>
        <!-- 每周一 9:00 -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/zhen/Documents/Codex/2026-06-07/deepseek/outputs/data/cron_weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zhen/Documents/Codex/2026-06-07/deepseek/outputs/data/cron_weekly.log</string>
</dict>
</plist>
